# dev/docker/supervisord-haproxy.conf
[program:haproxy]
command=/haproxy.sh
numprocs=1
autostart=true
autorestart=true
